# backend/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:app]
command=gunicorn main:app -w 4 -k uvicorn.workers.UvicornWorker -b 0.0.0.0:8000
directory=/app
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log

[program:celery]
command=celery -A main.celery worker --loglevel=info
directory=/app
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/celery.err.log
stdout_logfile=/var/log/celery.out.log

[program:flower]
command=celery -A main.celery flower --port=5555
directory=/app
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/flower.err.log
stdout_logfile=/var/log/flower.out.log

[program:backup]
command=python -m scripts.schedule_backup
directory=/app
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/backup.err.log
stdout_logfile=/var/log/backup.out.log 